# lleva_conexio.service
[Unit]
After=network.target

[Service]
User=root
ExecStart=/home/Administrador/llanca_examen_bd/lleva_conexio.sh

[Install]
WantedBy=default.target

